# dist/systemd/emru.service
[Unit]
After=network.target emru.socket
Requires=emru.socket

[Service]
ExecStart=/usr/bin/emrud
