# dist/systemd/emru.service
[Service]
ExecStartPre=/bin/rm -f /tmp/emru.sock
ExecStart=/usr/lib/emru/emrud
User=mohammad
Group=mohammad
Restart=on-failure

[Install]
WantedBy=multi-user.target
